# deploy/hoopci-gunicorn.service
# /etc/systemd/system/hoopci-gunicorn.service
[Unit]
Description=HoopCI — API Django (gunicorn)
After=network.target postgresql.service redis-server.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/hoopci/backend
Environment=DJANGO_SETTINGS_MODULE=hoopci.settings.production
ExecStart=/var/www/hoopci/backend/.venv/bin/gunicorn hoopci.wsgi:application -c gunicorn.conf.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
